# template/app.service.j2
[Unit]
Description=python application service

[Service]
ExecStart=/opt/application/run
Restart=on-abort

[Install]
WantedBy=multi-user.target
